# deploy/systemd/hintgrid-run.service
# Example systemd unit: one-shot pipeline run (install to /etc/systemd/system/).
# systemd does not start a second concurrent instance of this unit while ExecStart
# is still running (Type=oneshot). Adjust User, paths, and After= for your host.
#
# Install:
#   sudo cp deploy/systemd/hintgrid-run.service /etc/systemd/system/
#   sudo systemctl daemon-reload

[Unit]
Description=HintGrid batch pipeline (hintgrid run)
Documentation=file:///opt/hintgrid/README.ru.md
After=network-online.target
Wants=network-online.target
# Uncomment if PostgreSQL / Redis run on the same host (adjust unit names):
# After=postgresql.service redis-server.service

[Service]
Type=oneshot
User=hintgrid
Group=hintgrid
WorkingDirectory=/opt/hintgrid
EnvironmentFile=-/opt/hintgrid/.env
ExecStart=/opt/hintgrid/venv/bin/hintgrid run
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=4
StandardOutput=journal
StandardError=journal
SyslogIdentifier=hintgrid
TimeoutStartSec=infinity
TimeoutStopSec=300

[Install]
WantedBy=multi-user.target
