# deploy/systemd/hintgrid-run.service
# HintGrid batch pipeline (install to /etc/systemd/system/).
# Same host as Mastodon: PostgreSQL + Redis (systemd), Neo4j (Docker), mastodon-web + sidekiq.
# systemd does not start a second instance while ExecStart runs (Type=oneshot).
#
# Install:
#   sudo cp deploy/systemd/hintgrid-run.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#
# If PostgreSQL uses only a versioned unit (e.g. postgresql@16-main.service), replace
# postgresql.service in After=/Wants= (systemctl list-units 'postgresql*').
# If Redis unit is named redis.service, replace redis-server.service accordingly.

[Unit]
Description=HintGrid batch pipeline (hintgrid run)
Documentation=file:///opt/hintgrid/README.ru.md
After=network-online.target postgresql.service redis-server.service docker.service mastodon-web.service mastodon-sidekiq.service
Wants=network-online.target postgresql.service redis-server.service docker.service

[Service]
Type=oneshot
User=hintgrid
Group=hintgrid
WorkingDirectory=/opt/hintgrid
EnvironmentFile=-/opt/hintgrid/.env
ExecStart=/opt/hintgrid/venv/bin/hintgrid run
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=4
StandardOutput=journal
StandardError=journal
SyslogIdentifier=hintgrid
TimeoutStartSec=infinity
TimeoutStopSec=300

[Install]
WantedBy=multi-user.target
